# systemd/sigint-analyzer.service
[Unit]
Description=SIGINT Keystroke Analyzer
After=sigint-collector.service

[Service]
Type=exec
ExecStart=/usr/lib/sigint/sigint-analyzer
Restart=always
RestartSec=1

User=sigint-analyzer
Group=sigint

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateNetwork=true
RestrictNamespaces=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectHostname=true

RuntimeDirectory=sigint
RuntimeDirectoryMode=0750
StateDirectory=sigint
StateDirectoryMode=0700

ReadWritePaths=/var/lib/sigint

[Install]
WantedBy=multi-user.target
