# systemd/sigint-collector.service
[Unit]
Description=SIGINT Keystroke Collector
After=systemd-udevd.service
Wants=sigint-analyzer.service

[Service]
Type=exec
ExecStart=/usr/lib/sigint/sigint-collector
Restart=always
RestartSec=1

User=sigint-collector
Group=input
CapabilityBoundingSet=CAP_DAC_READ_SEARCH
AmbientCapabilities=CAP_DAC_READ_SEARCH
NoNewPrivileges=true

ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateNetwork=true
RestrictNamespaces=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectHostname=true

RuntimeDirectory=sigint
RuntimeDirectoryMode=0750

[Install]
WantedBy=multi-user.target
